# behind-repos/app-accessibility/espeakup/files/espeakup
#!/sbin/runscript
# Copyright 2010 Agilia Linux team/fat0troll
# Distributed under the terms of the GNU General Public License, v2 or later
# Special thanks goes to Gentoo

depend() {
	use logger acpid
	after coldplug dns nscd consolefont keymaps
}

start() {
	ebegin "Starting ESpeakup..."
	start-stop-daemon --start -q --exec /usr/bin/espeakup \
		--pidfile /var/run/espeakup.pid -b \
		-- ${ESPEAKUP_OPTS}
}

stop() {
	ebegin "Stopping ESpeakup..."
	start-stop-daemon --stop -q --pidfile /var/run/espeakup.pid 
}
